# templatefs/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:-/bin/sh
ttyS0::respawn:-/bin/sh
